# home/systemd/.config/systemd/user/rclone.timer
[Unit]
Description=timer for rclone

[Timer]
OnBootSec=2min

[Install]
WantedBy=timers.target
