# home/twitch-recorder/.config/systemd/user/halo.timer
[Unit]
Description=Run every minute

[Timer]
OnCalendar=*:*:0
Unit=halo.service

[Install]
WantedBy=timers.target

