# root/tunnel/etc/systemd/system/seedbox-tunnel.service
[Unit]
Description=Runs multiple SSH tunnels with autossh
After=network.target

[Service]
User=qwerty
Group=qwerty
Environment="AUTOSSH_GATETIME=0"
ExecStart=/usr/bin/autossh -M 0 -N -D 50203 -o "ControlMaster=no" -o "ExitOnForwardFailure=yes" -o "ServerAliveInterval=30" -o "ServerAliveCountMax=3" -l qwertyspace gegenees

[Install]
WantedBy=multi-user.target
