# systemd/.config/systemd/user/slideshow.timer
[Unit]
Description=Runs the slideshow service

[Timer]
OnBootSec=0s
OnUnitActiveSec=300s

[Install]
WantedBy=timers.target
